# noninteractive/.crontab
#
# $Id$
#
00 02 * * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/VbRescueCouncil/noninteractive/milestone_sweep.aspx
15 02 * * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/VbRescueCouncil/noninteractive/daily_tasks.aspx
